# python/resources/systemd-user/dynamic_power_command.service
[Unit]
Description=Dynamic‑Power tray UI
After=graphical-session.target
Requires=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/dynamic_power_command
Restart=on-failure

[Install]
WantedBy=graphical-session.target
